# config/deployment/systemd/pynomaly-dashboard.service
[Unit]
Description=Pynomaly Monitoring Dashboard
After=network.target
Requires=network.target

[Service]
Type=simple
User=pynomaly
Group=pynomaly
WorkingDirectory=/mnt/c/Users/andre/Pynomaly
Environment=PYTHONPATH=/mnt/c/Users/andre/Pynomaly/src
ExecStart=/usr/bin/python3 -m pynomaly.infrastructure.monitoring.dashboard
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
