# src/temporary/production_config/pynomaly.service

[Unit]
Description=Pynomaly Production Service
After=network.target postgresql.service redis.service

[Service]
Type=forking
User=pynomaly
Group=pynomaly
WorkingDirectory=/app
Environment=PYNOMALY_ENV=production
ExecStart=/usr/local/bin/gunicorn -c gunicorn.conf.py pynomaly.presentation.api.app:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure

[Install]
WantedBy=multi-user.target
